# netnoise.service
[Unit]
Description=netnoise Service
Documentation=https://github.com/madebyjake/netnoise
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/netnoise
ExecStart=/opt/netnoise/netnoise.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netnoise

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/netnoise /tmp

# Resource limits
MemoryMax=256M
CPUQuota=50%

# Restart settings
Restart=on-failure
RestartSec=30
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
